# deploy/linux/pc-inventory-agent.service
# systemd-Unit für den Roster-Agent.
# Installation:
#   sudo cp agent-linux-amd64 /usr/local/bin/roster-agent
#   sudo mkdir -p /etc/roster && sudo cp agent.yaml.example /etc/roster/agent.yaml
#   # agent.yaml anpassen (server_url + enrollment_token)
#   sudo cp roster-agent.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now roster-agent
#
# Alternativ registriert der Agent den Dienst selbst:
#   sudo /usr/local/bin/roster-agent -config /etc/roster/agent.yaml install

[Unit]
Description=Roster Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/roster-agent -config /etc/roster/agent.yaml run
Restart=on-failure
RestartSec=10
# Eingeschränkte Rechte; für die DMI-Seriennummer ggf. anpassen.
DynamicUser=yes
StateDirectory=roster

[Install]
WantedBy=multi-user.target
