# deploy/macos/com.pc-inventory.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd-Daemon für den Roster-Agent (macOS).
  Installation:
    sudo cp agent-darwin-arm64 /usr/local/bin/roster-agent
    sudo mkdir -p /etc/roster && sudo cp agent.yaml.example /etc/roster/agent.yaml
    # agent.yaml anpassen (server_url + enrollment_token)
    sudo cp com.roster.agent.plist /Library/LaunchDaemons/
    sudo launchctl load /Library/LaunchDaemons/com.roster.agent.plist

  Alternativ registriert der Agent den Dienst selbst:
    sudo /usr/local/bin/roster-agent -config /etc/roster/agent.yaml install
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.roster.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/roster-agent</string>
        <string>-config</string>
        <string>/etc/roster/agent.yaml</string>
        <string>run</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/var/log/roster-agent.log</string>
    <key>StandardOutPath</key>
    <string>/var/log/roster-agent.log</string>
</dict>
</plist>
